# packaging/stampede/main/src/main/dist/systemd/torodb-stampede.service.sample
[Unit]
Description=ToroDB
After=syslog.target network.target postgresql.target

[Service]
Type=simple
User=torodb
ExecStart=/usr/bin/@{assembler.name}
PrivateTmp=true
LimitNOFILE=64000
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target